# initramfs/inittab
::sysinit:/etc/init.d/rcS
::restart:/init
::ctrlaltdel:/sbin/reboot